# supervisord.conf
[supervisord]
nodaemon=true

[program:cassandra]
command=cassandra -f

[program:sparkWorker]
command=/start-worker.sh
