# supervisor.conf/supervisord-master.conf
[supervisord]
nodaemon=true

[program:sparkMaster]
command=/start-master.sh
