# supervisor.conf/supervisord-cass.conf
[supervisord]
nodaemon=true

[program:cassandra]
command=cassandra -f -R
